# service/del_image.timer
[Unit]
Description=Run Delete Images Script Every Hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
